# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#NEXSEC add for ipsec performance tunning
/etc/ucpe_tuning.sh
exit 0
